# roles/wallpaper/templates/lib/systemd/system/ofi-wallpaper.service.j2
[Unit]
Description={{ service_name }}

[Service]
Environment=DISPLAY=:0
ExecStart=/bin/bash -c "/{{ service_dir }}/{{ service_name }}.sh"

[Install]
WantedBy=basic.target

